# etc/supervisord.d/zookeeper_update.ini
[program:zookeeper-update]
command=consul-template -consul localhost:8500 -wait=60s -template "/etc/consul-templates/zoo.cfg.ctmpl:/opt/zookeeper/conf/zoo.cfg:supervisorctl restart zookeeper"
redirect_stderr=true
autostart=false
stdout_logfile=/var/log/supervisor/%(program_name)s.log
