# etc/supervisord.d/zookeeper_update.ini
[program:zookeeper-update]
command=/opt/qnib/zookeeper/bin/start_zookeeper-update.sh
redirect_stderr=true
exitcodecs=0
startsecs=0
autostart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
